# systemd/user/rclone-proton.timer
[Unit]
Description=Run rclone-proton periodically

[Timer]
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
